# files/upload-snapshots.service
[Unit]
Description=Snapshot service

[Service]
ExecStart=/home/ubuntu/uploadsnapshots.sh

[Install]
WantedBy=multi-user.target

